# powertop.service
[Unit]
Description=powertop tunings

[Service]
Type=oneshot
ExecStart=/home/carlbernal/scripts/pt-at

[Install]
WantedBy=multi-user.target
